# deploy/host/finpipe.service
[Unit]
Description=finpipe docker compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=@@INSTALL_DIR@@
User=@@USER@@
EnvironmentFile=@@INSTALL_DIR@@/deploy/secrets/.env
ExecStartPre=/usr/bin/git -C @@INSTALL_DIR@@ pull --ff-only
ExecStart=/usr/bin/docker compose -f deploy/compose/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f deploy/compose/docker-compose.yml down

[Install]
WantedBy=multi-user.target
